# etc/udev/rules.d/95-monitors.rules
ACTION=="change", RUN+="/home/zaq/scripts/disp-auto"
